# scripts/powertop.service
[Unit]
Description=Powertop tunings

[Service]
ExecStart=

[Install]
WantedBy=multi-user.target